# crontab
# Every Hour
0 * * * * /usr/bin/python /home/tanmaythakur2/github_committer/github_committer.py

# Every Midnight
0 0 * * * /usr/bin/python /home/tanmaythakur2/github_committer/randomness.py

